# implementation/systemd/multi-agent-pipeline-watchdog.timer
[Unit]
Description=Run Multi-Agent Research Pipeline Watchdog every five minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true
Unit=multi-agent-pipeline-watchdog.service

[Install]
WantedBy=timers.target
