# implementation/systemd/multi-agent-pipeline-notify.timer
[Unit]
Description=Flush research pipeline notifications

[Timer]
OnBootSec=45s
OnUnitInactiveSec=30s
AccuracySec=1s
Persistent=true
Unit=multi-agent-pipeline-notify.service

[Install]
WantedBy=timers.target
